# alt/deploy_files/etc/systemd/system/ilxalt.socket
[Unit]
Description=interlex alt gunicorn socket

[Socket]
ListenStream=/run/ilxalt/socket

[Install]
WantedBy=sockets.target
